# .config/systemd/user/sommelier-x@0.service.d/override.conf
[Service]
Environment="SOMMELIER_ACCELERATORS=Super_L,<Alt>tab"
